# rpm/evebox-agent.service
[Unit]
Description=EveBox Agent

[Service]
ExecStart=/usr/bin/evebox agent

[Install]
WantedBy=multi-user.target
